# tutorial/old_code/train_image_2024/config/supervisor/ssh.conf
[program:ssh]
command = /bin/bash /opt/config/supervisor/ssh.sh
autostart = true
autorestart = true
stopasgroup = true
killasgroup = true